# scripts/linux/slipstream.service
[Unit]
Description=SlipStream Media Management
Documentation=https://github.com/slipstream/slipstream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/slipstream
WorkingDirectory=%h/.config/slipstream
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
